# sdbuild/packages/custom-settings/tcf.agent.service
[Unit]
Description=TCF agent daemon

[Service]
Type=forking
ExecStart=/usr/local/sbin/tcf-agent -d

[Install]
WantedBy=basic.target
